# source/profile.d/snd_usb_card.sh

systemctl --user restart 'snd_usb_card@*'  --all
